# iRedMail/samples/policyd-cleanup.cron
#!/bin/sh

# Author: Zhang Huangbin (michaelbibby <at> gmail.com)

#
# Sample policyd-cleanup crontab script. It should be localted at:
#   /etc/cron.daily/policyd-cleanup
#
# Shipped within iRedMail project:
#   * http://iRedMail.googlecode.com/
#

/usr/sbin/policyd-cleanup -c /etc/policyd.conf
/usr/sbin/policyd-cleanup -c /etc/policyd_sender_throttle.conf
